# dynamo/MCRLinux/sys/ros1/glnxa64/ros1/etc/catkin/profile.d/1.ros_package_path.sh
version https://git-lfs.github.com/spec/v1
oid sha256:a4e5e3705eee93522c1d5a7186a0f7995d6c9167177e4362f627e75ecc1462bd
size 1045
