# dynamo/MCRLinux/sys/ros1/glnxa64/ros1/etc/catkin/profile.d/1.ros_version.sh
version https://git-lfs.github.com/spec/v1
oid sha256:bf8ab9eb545afd016930d9c6f8f7ce713b356ab3ecf3f48c6efeca23676af570
size 85
