# dynamo/MCRLinux/sys/ros1/glnxa64/ros1/etc/catkin/profile.d/1.ros_distro.sh
version https://git-lfs.github.com/spec/v1
oid sha256:2fcb7b1b54373521c2ad459f9c45a4133716b4eda09e1780db1fe335c0f1c3c0
size 290
