# dynamo/MCRLinux/sys/ros1/glnxa64/ros1/etc/catkin/profile.d/15.rosbash.fish
version https://git-lfs.github.com/spec/v1
oid sha256:35f5c282328dec54ec38ebde0e2e364598478adb5533a7eb8b0958be2f1c864d
size 206
