# dynamo/MCRLinux/sys/ros1/glnxa64/ros1/etc/catkin/profile.d/15.rosbash.tcsh
version https://git-lfs.github.com/spec/v1
oid sha256:1273e3e0a2f948e55e8bdbc10b678a85b574f6eea852c1776a42a0e2133e1b60
size 206
